# lebackup.cron
#!/bin/bash
/opt/lebackup/lebackup-all.sh
